# workspace/crontab/laradock
# * * * * * laradock /usr/bin/php /var/www/mall/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/mall/artisan schedule:run >> /var/log/crontest.log
